# reboot.cron
* * * * * sudo /home/pi/orobot-firmware/reboot.sh
@reboot /home/pi/orobot-firmware/reboot.sh
