# reboot.cron
* * * * * /home/pi/projects/orobot-firmware/reboot.sh
