# autostart/display.desktop
[Desktop Entry]
Name=OROBOTLogs
Exec=/home/pi/orobot-firmware/init.d/orobot-boot
